# production/backup.service
[Unit]
Description=Backup

[Service]
Type=oneshot
ExecStart=/home/loconomics/backup.sh
Nice=19
